# s6-rc.d/init-tacplus/run
#!/command/with-contenv sh
/usr/local/sbin/tac_plus-ng -P $TACPLUS_CFG_FILE
