# s6-rc.d/svc-tacplus/run
#!/command/with-contenv sh
/usr/local/sbin/tac_plus-ng -f $TACPLUS_CFG_FILE
